# worker/deepseek-worker.service
[Unit]
Description=DeepSeek Orchestration Worker Daemon
After=network.target

[Service]
Type=simple
User=deepseek
Group=deepseek
WorkingDirectory=/opt/deepseek-worker

# Environment variables (AWS Role, Azure Vault URL, LM Studio URL)
EnvironmentFile=/etc/deepseek/.env

# Start the Python daemon
ExecStart=/opt/deepseek-worker/venv/bin/python /opt/deepseek-worker/deepseek_daemon.py

# Crash Recovery: Auto-restart on failure
Restart=always
RestartSec=10

# Security Hardening
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=yes

[Install]
WantedBy=multi-user.target
